# tomcat/supervisord.conf
[supervisord]
nodaemon=true

[program:Tomcat7]
command=service tomcat7 start --privileged=true && tail -f /var/log/tomcat7/catalina.out